# kube/research-ai-ai.service
[Unit]
Description=Research-AI AI Service (llama-server + FastAPI)
After=network.target

[Service]
Type=simple
# Change to the user that owns the Research-AI repo checkout
User=lvandee
# All paths are configured here — change these for different deployments
EnvironmentFile=-/opt/research-ai/ai.env
ExecStart=%h/Research-AI/ai/scripts/start-ai.sh
ExecStop=%h/Research-AI/ai/scripts/stop-ai.sh
Restart=on-failure
RestartSec=10
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
